# cookbooks/cmdr-defaults/files/default/lirc.conf
description "lircd"

start on (started iguanaIR)
stop on runlevel [06]
expect fork
respawn

pre-start script
		mkdir -p /var/run/lirc
end script


exec lircd -H iguanaIR -d /dev/iguanaIR/0
