# cookbooks/roomtrol-randoms/files/default/iguanaIR.conf
description "iguanaIR"

start on (filesystem and started dbus)
stop on runlevel [06]
expect fork
respawn

exec igdaemon
