# cookbooks/roomtrol-randoms/files/default/roomtrol-daemon.conf
description "Roomtrol daemon script"

respawn

start on (started rc-sysinit)
stop on runlevel [06]

exec bootup_roomtrol
