# cookbooks/chef-client/templates/default/com.opscode.chef-client.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

<dict>
  <key>Label</key>
  <string>com.opscode.chef-client</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/chef-client</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>UserName</key>
  <string>root</string>
  <key>StandardOutPath</key>
  <string><%= node["chef_client"]["log_dir"] %>/client.log</string>
  <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>00</integer>
        <key>Hour</key>
        <integer>3</integer>
    </dict>
</dict>

</plist>
